# pm-reportemail.service
[Unit]
Description=Power monitor email script

[Service]
Type=simple
ExecStart=/usr/bin/pm-reportemail
