# pm-reportemail.service
[Unit]
Description=Power monitor email script

[Service]
Type=simple
ExecStart=/usr/bin/python /home/powermonitor/pm-reportemail.py
